# docker/services/memcached/run
#!/bin/sh
/sbin/setuser memcache /usr/bin/memcached

