# docker/services/alpharpc-clienthandler/run
#!/bin/sh
exec /sbin/setuser alpharpc /var/www/bin/alpharpc-client-handler
